# templates/default/postgresql.service.erb
.include /usr/lib/systemd/system/postgresql.service
[Service]
Environment=
Environment=PGPORT=<%= node['postgresql']['config']['port'] %>
Environment=PGDATA=<%= node['postgresql']['config']['data_directory'] %>
