# templates/default/postgresql.service.erb
[Service]
Environment=
Environment=PGPORT=<%= node['postgresql']['config']['port'] %>
Environment=PGDATA=<%= node['postgresql']['config']['data_directory'] %>
