# templates/default/postgresql.service.erb
[Service]
<% if node['postgresql']['use_pgdg_packages'] %>
.include /usr/lib/systemd/system/postgresql-<%= node['postgresql']['version'] %>.service
<% else %>
.include /usr/lib/systemd/system/postgresql.service
<% end %>

Environment=
Environment=PGPORT=<%= node['postgresql']['config']['port'] %>
Environment=PGDATA=<%= node['postgresql']['config']['data_directory'] %>
