# templates/default/postgresql.service.erb
.include /lib/systemd/system/postgresql-<%= node['postgresql']['version'] %>.service
[Service]
Environment=PGDATA=<%= node['postgresql']['dir'] %>
