# templates/default/postgresql.service.erb
[Service]
.include /usr/lib/systemd/system/<%= @svc_name %>.service

Environment=
Environment=PGPORT=<%= @port %>
Environment=PGDATA=<%= @data_dir %>
