# templates/default/postgresql.service.erb
.include /lib/systemd/system/<%= @path %>.service
[Service]
Environment=PGDATA=<%= node['postgresql']['dir'] %>
